# surface/macosx/field/ressources/fuse4x.bundle/Library/LaunchAgents/org.fuse4x.autoupdater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.fuse4x.autoupdater</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Filesystems/fuse4x.fs/Contents/Executables/autoupdater</string>
	</array>
	<key>StartInterval</key>
	<integer>864000</integer>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
